# scripts/systemd/metiqd.service
[Unit]
Description=metiq daemon — Nostr-native AI agent
Documentation=https://github.com/metiq/metiq
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Load configuration from ~/.metiq/.env
EnvironmentFile=-%h/.metiq/.env

# The binary is expected to be on PATH or at ~/.local/bin/metiqd
ExecStart=%h/.local/bin/metiqd

# Restart policy
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Hardening (user service — no root needed)
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/.metiq %h/.local/share/swarmstr /tmp

[Install]
WantedBy=default.target
